# ops/ci/uruquim-ci.timer
[Unit]
Description=Run the Uruquim verification gate every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=uruquim-ci.service

[Install]
WantedBy=timers.target
